# crontab
* * * * * curl --connect-timeout 2 --silent --fail --out /dev/null http://tales.squirr.us/tick
